# heaven/debian/cthulhu-heaven.service
[Unit]
Description=Cthulhu - Heaven

[Service]
ExecStart=/usr/bin/cthulhu-heaven -c /etc/cthulhu/heaven.toml

[Install]
WantedBy=multi-user.target
